# docker/crontab
# Harmonix monitoring cron jobs
# Position snapshots (every 5 min)
*/5 * * * * cd /app && env $(cat /etc/environment | xargs) python scripts/pull_positions_v3.py >> /app/logs/pull_positions.log 2>&1

# Position prices (every 5 min, offset by 2 min)
2-57/5 * * * * cd /app && env $(cat /etc/environment | xargs) python scripts/pull_position_prices.py >> /app/logs/pull_position_prices.log 2>&1

# Market data (hourly at :20)
20 * * * * cd /app && env $(cat /etc/environment | xargs) python scripts/pull_hyperliquid_v3.py >> /app/logs/pull_hyperliquid.log 2>&1

# Computation pipeline (hourly at :30)
30 * * * * cd /app && env $(cat /etc/environment | xargs) python scripts/pipeline_hourly.py >> /app/logs/pipeline_hourly.log 2>&1

# Empty line required by cron
